# deploy/launchd/com.aiopscodex.host-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aiopscodex.host-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/lizhongxuan/Desktop/aiops-codex/.data/bin/host-agent</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AIOPS_SERVER_GRPC_ADDR</key>
    <string>127.0.0.1:18090</string>
    <key>AIOPS_AGENT_BOOTSTRAP_TOKEN</key>
    <string>change-me</string>
    <key>AIOPS_AGENT_HOST_ID</key>
    <string>agent-local</string>
    <key>AIOPS_AGENT_HOSTNAME</key>
    <string>agent-local</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/lizhongxuan/Desktop/aiops-codex</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/lizhongxuan/Desktop/aiops-codex/.data/launchd-host-agent.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/lizhongxuan/Desktop/aiops-codex/.data/launchd-host-agent.err.log</string>
</dict>
</plist>
